# etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:bird]
command=/usr/local/sbin/bird -c /etc/bird/bird.conf

[program:birdwatcher]
command=/usr/bin/birdwatcher -config /etc/birdwatcher/birdwatcher.conf